# etc/profile.d/fastfetch2077.sh
#!/usr/bin/env bash
if [ -x "$HOME/bin/fastfetch2077.sh" ]; then
  $HOME/bin/fastfetch2077.sh
fi
